# telegram-bot.service
[Unit]
Description=Telegram Bot Admin Panel
After=network.target
Wants=network.target

[Service]
Type=simple
User=telegram_bot_admin
Group=telegram_bot_admin
WorkingDirectory=/home/telegram_bot_admin
Environment=PATH=/home/telegram_bot_admin/venv/bin
Environment=PYTHONPATH=/home/telegram_bot_admin
ExecStart=/home/telegram_bot_admin/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Автоматический перезапуск при сбоях
RestartPreventExitStatus=0
RestartForceExitStatus=1

# Ограничения ресурсов
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
